# etc/cron.d/sssync
57 * * * *    root    /home/toveri/SSAuthenticator/scripts/sssync
